# root/etc/s6-overlay/s6-rc.d/syslog-ng-init/up
with-contenv
exec /etc/s6-overlay/s6-rc.d/syslog-ng-init/syslog-ng-init.sh
