# root/etc/s6-overlay/s6-rc.d/syslog-ng/run
#!/usr/bin/execlineb -P
/usr/sbin/syslog-ng -F
